# conf.d/monitor_tasks.conf
[program:monitor_tasks]
command=python3 manage.py monitor_tasks
autostart=true
autorestart=true
stderr_logfile=/var/log/monitor_tasks.err.log
stdout_logfile=/var/log/monitor_tasks.out.log
