# conf.d/process_tasks.conf
[program:process_tasks]
command=python3 manage.py process_tasks
autostart=true
autorestart=true
stderr_logfile=/var/log/process_tasks.err.log
stdout_logfile=/var/log/process_tasks.out.log

